# var/spool/cron/crontabs/root
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/var/scripts/supervision/92.103.197.250:65503 installed on Fri May 19 11:05:01 2017)
# (Cron version -- $Id: crontab.c,v 2.13 1994/01/17 03:20:37 vixie Exp $)
*/5 * * * * /var/scripts/script_supervision.sh
*/5 * * * * /var/scripts/script_gestion_logs_1.sh
*/5 * * * * /var/scripts/script_tueur.sh
*/5 * * * * /var/scripts/supervision/92.103.197.250:51461.sh
*/5 * * * * /var/scripts/supervision/92.103.197.250:51461.sh
*/5 * * * * /var/scripts/supervision/92.103.197.250:65503.sh
